# shop1/meituan-browser-control.service
[Unit]
Description=Meituan Waimai browser control
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/meituan-reply-bot
Environment=PYTHONUNBUFFERED=1
Environment=DISPLAY=:99
ExecStartPre=/home/ubuntu/meituan-reply-bot/kill_chromium.sh
ExecStartPre=/bin/sleep 2
ExecStart=/home/ubuntu/meituan-reply-bot/remote_browser.sh
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
