# etc/meituan-mem-watch.service
[Unit]
Description=Meituan bot memory watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/meituan-mem-watch.sh
