# etc/meituan-cookie-watch@.timer
[Unit]
Description=Meituan cookie auto-refresh timer (%i)

[Timer]
OnBootSec=5min
OnUnitActiveSec=20h
AccuracySec=5min
Unit=meituan-cookie-watch@%i.service
Persistent=true

[Install]
WantedBy=timers.target
