# etc/meituan-cookie-watch@.service
[Unit]
Description=Meituan cookie auto-refresh watcher (%i)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/meituan-cookie-watch.sh %i
Nice=10
