# etc/meituan-promo-scheduler.service
[Unit]
Description=Meituan Waimai promotion scheduler
After=meituan-browser-control.service
Wants=meituan-browser-control.service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/meituan-reply-bot
ExecStart=/home/ubuntu/meituan-reply-bot/.venv/bin/python -u /home/ubuntu/meituan-reply-bot/promo_scheduler.py --config /home/ubuntu/meituan-reply-bot/config.yaml
Restart=always
RestartSec=10
User=ubuntu

[Install]
WantedBy=multi-user.target
